# deploy/forecast_runner/systemd/forecast-ops-backend.service.in
[Unit]
Description=Forecast Ops Codex execution backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@RUNNER_USER@
Group=@RUNNER_USER@
WorkingDirectory=@RUNNER_ROOT@
EnvironmentFile=@BACKEND_ENV_FILE@
Environment=HOME=@RUNNER_ROOT@/home
Environment=CODEX_HOME=@RUNNER_ROOT@/codex-home
Environment=GH_CONFIG_DIR=@RUNNER_ROOT@/github
Environment=PATH=@RUNNER_ROOT@/.runtime/codex/node_modules/.bin:@RUNNER_ROOT@/.runtime/tools/node_modules/.bin:@RUNNER_ROOT@/.runtime/bootstrap/bin:@RUNNER_ROOT@/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=@RUNNER_ROOT@/backend/run.sh
Restart=on-failure
RestartSec=5s
KillMode=control-group
TimeoutStopSec=90s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=@RUNNER_ROOT@
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
