# deploy/forecast_runner/systemd/forecast-replica-backup.timer.in
[Unit]
Description=Daily Forecast Ops replica bundle

[Timer]
OnCalendar=*-*-* 19:00:00 UTC
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
